# .config/supervisor/uwsgi.conf
[program:uwsgi]
command = /root/.pyenv/versions/used-book-store/bin//uwsgi --ini /srv/used-book-store/.config/uwsgi/uwsgi-app.ini